# templates/restic-forget.timer.j2
[Unit]
Description=Restic backup forget

[Timer]
OnCalendar={{ restic_backup_schedule }}
RandomizedDelaySec={{ restic_systemd_randomised_delay }}
Persistent=true

[Install]
WantedBy=timers.target
